# task.cron
0 * * * * /bin/sh /home/azureuser/LuisE2E/runAndSendMail.sh
