# task.cron
30 3 * * * /bin/sh /home/azureuser/LuisE2E/runAndSendMail.sh >> /home/azureuser/LuisE2E/e2etest.log 2>&1
